# backend/sigma_rules/sigma_rule_008_suspicious_scheduled_task.yml
title: Suspicious Scheduled Task Creation
id: sigma-008
status: experimental
description: Detects creation of suspicious scheduled tasks for persistence
references:
    - https://attack.mitre.org/techniques/T1053/005/
tags:
    - attack.execution
    - attack.persistence
    - attack.t1053.005
logsource:
    product: windows
    service: task_scheduler
detection:
    selection:
        - 'schtasks /create'
        - 'powershell.exe'
        - 'cmd.exe /c'
    condition: selection
falsepositives:
    - Legitimate automation scripts
    - System maintenance tasks
level: medium

